# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:postgresql]
command=su -s /bin/bash postgres -c "/usr/lib/postgresql/16/bin/postgres -D %(ENV_PGDATA)s"
autostart=true
autorestart=true
priority=10
stdout_logfile=/var/log/supervisor/postgresql.log
stderr_logfile=/var/log/supervisor/postgresql.log
redirect_stderr=true
environment=PGDATA="%(ENV_PGDATA)s"

[program:node_exporter]
command=/opt/node_exporter/node_exporter
    --web.listen-address=0.0.0.0:9100
    --collector.filesystem.mount-points-exclude=^/(sys|proc|dev|run)($|/)
autostart=true
autorestart=true
priority=12
stdout_logfile=/var/log/supervisor/node_exporter.log
stderr_logfile=/var/log/supervisor/node_exporter.log
redirect_stderr=true

[program:prometheus]
command=/opt/prometheus/prometheus
    --config.file=/etc/prometheus/prometheus.yml
    --storage.tsdb.path=/var/lib/prometheus
    --storage.tsdb.retention.time=30d
    --storage.tsdb.wal-compression
    --web.listen-address=0.0.0.0:9090
autostart=true
autorestart=true
priority=15
startsecs=5
startretries=10
stdout_logfile=/var/log/supervisor/prometheus.log
stderr_logfile=/var/log/supervisor/prometheus.log
redirect_stderr=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=19
stdout_logfile=/var/log/supervisor/nginx.log
stderr_logfile=/var/log/supervisor/nginx.log
redirect_stderr=true

[program:spring]
command=/start-spring.sh
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/supervisor/spring.log
stderr_logfile=/var/log/supervisor/spring.log
redirect_stderr=true
startretries=5
startsecs=15

